# src/99-betterusbcopy.rules
ACTION=="add", KERNEL=="sd[a-z][0-9]", RUN+="/bin/bash -c '/usr/local/bin/betterusbcopy &'"
